# curriculum-server.service
[Unit]
Description=AKA tuned-resonator curriculum-barnacle server service
After=curriculum-syslog.service
Requires=curriculum-syslog.service

[Service]
Type=simple
WorkingDirectory=/root/tuned-resonator/server/
ExecStart=/usr/bin/python server.py
Restart=on-abort
RestartSec=5
KillMode=process

[Install]
WantedBy=multi-user.target
